# systemd/server-pulse.timer
[Unit]
Description=Run server-pulse health checks every 5 minutes
Documentation=https://github.com/volkv/server-pulse

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
AccuracySec=10s
Unit=server-pulse.service
Persistent=true

[Install]
WantedBy=timers.target
